# etc/cron.d/diskcleaner
0 20 * * * root /usr/local/bin/diskcleaner.sh
